# custom_systemd_services/.config/systemd/user/obsidian-auto-commit.service
[Unit]
Description=Automated hourly commit and push for Obsidian vault

[Service]
ExecStart=/home/digs/bin/obsidian-auto-commit.sh  # Full path to your script
WorkingDirectory=/home/digs/Vaults/diego  # Optional, but ensures correct dir
