# custom_systemd_services/.config/systemd/user/obsidian-auto-commit.timer
[Unit]
Description=Run Obsidian auto-commit every hour on the hour

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=false
Unit=obsidian-auto-commit.service

[Install]
WantedBy=timers.target
